# docker-import/systemd/opencode.service
[Unit]
Description=OpenCode API Server
After=network.target

[Service]
Type=simple
User=will
WorkingDirectory=/home/will
ExecStart=/home/will/.opencode/bin/opencode serve --hostname 127.0.0.1 --port 3456
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
